# src/crontab
*/10 * * * * python3 /usr/local/price_updater/update_prices.py -c /usr/local/price_updater/config/config.json
